# units/rclone-sync-obsidian.timer
[Unit]
Description=rclone hourly backup - Obsidian

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
